# certbot/files/certbot.service
[Unit]
Description=Let's Encrypt renewal
[Service]
Type=oneshot
ExecStart=/opt/certbot/get-certs.sh